# server/lobster-analytics.service
[Unit]
Description=Lobster Farm Analytics Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/lobster-farm/server
ExecStart=/usr/bin/python3 /www/wwwroot/lobster-farm/server/analytics_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
